# deploy/systemd/bn-square-agent-mcp.service
[Unit]
Description=BN Square Agent self-hosted MCP publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bn-square
Group=bn-square
WorkingDirectory=/opt/bn-square-agent
EnvironmentFile=/etc/bn-square-agent/env
ExecStart=/opt/bn-square-agent/.venv/bin/python -B run.py serve-mcp --host 127.0.0.1 --port 8788
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bn-square-agent/data

[Install]
WantedBy=multi-user.target
